# ansible/roles/apollo-portal/files/apollo_portal.socket
[Unit]
Description=Apollo portal website socket
[Socket]
ListenStream=/run/apollo_portal.sock
[Install]
WantedBy=sockets.target
